# display-server.service
[Unit]
Description=AI E-Ink Display Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=austingibb
WorkingDirectory=/home/austingibb/ai_desk_agent
ExecStart=/home/austingibb/ai_desk_agent/venv/bin/python display_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
